# docker/etc/s6-overlay/s6-rc.d/svc-adonis/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/adonis || exit

exec s6-setuidgid root nodejs bin/server.js
