# docker/etc/s6-overlay/s6-rc.d/svc-express/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/express || exit

exec s6-setuidgid root ndoejs index.js
